# sv/adbd/finish
#!/system/bin/sh
stop adbd
setprop service.adb.tcp.port -1
settings put global adb_enabled 0
